# systemd/blockchain-monitor-mainnet.service
[Unit]
Description=Studio Blockchain Monitoring Service
After=geth-studio-mainnet.service
Wants=geth-studio-mainnet.service

[Service]
Type=simple
User=root
ExecStart=/bin/bash ~/studio-mainnet/node/scripts/monitor_blockchain.sh
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=blockchain-monitor-mainnet

# Hardening measures
PrivateTmp=true
ProtectSystem=full
NoNewPrivileges=true
ProtectHome=false

[Install]
WantedBy=multi-user.target
