# systemd/geth-studio.service
[Unit]
Description=Studio Blockchain Mainnet Node
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
ExecStart=/bin/bash /root/studio-mainnet/node/start-node.sh
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=geth-studio
LimitNOFILE=65536
LimitNPROC=65536

# Ensure data integrity during shutdown
KillSignal=SIGINT
TimeoutStopSec=300

# Hardening measures
PrivateTmp=true
ProtectSystem=full
NoNewPrivileges=true
ProtectHome=false

[Install]
WantedBy=multi-user.target
