# config/blockchain-monitor.service
[Unit]
Description=Studio Blockchain Monitor
After=network.target geth-studio.service

[Service]
Type=simple
User=root
ExecStart=/bin/bash /root/enhanced_monitor_blockchain.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
